# Services/after-wakeup.service
[Unit]
Description=Bruno Fontes script after wakeup
After=suspend.target
OnFailure=logErrors@%n.service

[Service]
Type=simple
ExecStart=/usr/bin/systemctl restart openvpn-client@privatvpn.service

[Install]
WantedBy=suspend.target
