# Services/backup-contabo-files.service
[Unit]
Description=Backup contabo files once per week

[Service]
ExecStart=/home/bruno/Backups/Contabo/backup-files.sh
