# Services/backup-contabo-db.service
[Unit]
Description=Backup contabo db daily

[Service]
ExecStart=/home/bruno/Backups/Contabo/backup-db.sh
