# Services/backup-contabo-db.timer
[Unit]
Description=Backup contabo db daily

[Timer]
Persistent=true
OnCalendar=*-*-* 10:00:00
WakeSystem=true

[Install]
WantedBy=timers.target
