# Services/backup-contabo.service
[Unit]
Description=Backup contabo daily

[Service]
ExecStart=/home/bruno/Backups/Contabo/backup-contabo.sh
